# deb-control/prerm
#!/bin/sh

systemctl stop piaq
systemctl disable piaq
echo "Removed piaq service"
